# deploy/recipes/solo/supervisord.conf
; DataHelix solo recipe — process supervision (proposal decision 1.5).
; Three processes, one container: the data plane, the web front, and the
; Modeler's git relay. Each restarts independently; all log to the container's
; stdout/stderr with a program prefix.

[program:mosaic]
command=%(ENV_MOSAIC_BIN)s serve --graphql --host 127.0.0.1 --port 8001
directory=/project
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cors-proxy]
command=/usr/local/bin/node /srv/cors-proxy/server.js
environment=PORT="9999",ALLOW_ORIGIN="*"
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
