# deploy/recipes/demo/supervisord.conf
; DataHelix demo recipe — process supervision. Two processes: the data plane
; (serving the baked dataset) and the web front. No cors-proxy (no Modeler).

[program:mosaic]
command=%(ENV_MOSAIC_BIN)s serve --graphql --config %(ENV_MOSAIC_CONFIG)s --host 127.0.0.1 --port 8001
directory=/project
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
